# dh-ddns-updater.service
[Unit]
Description=Dreamhost Dynamic DNS Updater
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=dh-ddns-updater
Group=dh-ddns-updater
ExecStart=/usr/local/bin/dh-ddns-updater
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/dh-ddns-updater
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target